# build/openrwc-amd64/usr/lib/openrwc/openrwc.service
[Unit]
Description=Reddit Wallpaper Changer for GNU/Linux
