# build/openrwc-amd64/usr/lib/openrwc/openrwc.service
[Unit]
Description=Reddit Wallpaper Changer for GNU/Linux
After=display-manager.service
StartLimitIntervalSec=0

[Service]
Environment="DISPLAY=:0.0"
Environment="XAUTHORITY=%h/.Xauthority"
Type=simple
ExecStart=/usr/bin/sh -c "DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus /usr/bin/openrwc"
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
