# build/openrwc-amd64/etc/systemd/system/openrwc@.service
[Unit]
Description=Reddit Wallpaper Changer for GNU/Linux
Requires=display-manager.service
After=display-manager.service
StartLimitIntervalSec=0

[Service]
Environment="DISPLAY=:0"
Type=simple
ExecStart=/usr/bin/openrwc
Restart=always
RestartSec=5
User=%i

[Install]
WantedBy=graphical.target
