# build/openrwc-amd64/DEBIAN/postinst
#!/bin/env bash

mkdir -p ~/.config/systemd/user
cp /usr/lib/openrwc/openrwc.service ~/.config/systemd/user/

chown -R "$USER":"$USER" ~/.config/systemd/

systemctl daemon-reload

echo "
OpenRWC Installed 🧞

Run the following commands to begin
"

echo -e "\tsystemctl --user enable openrwc.service"
echo -e "\tsystemctl --user start openrwc.service"

echo -e "\tjournalctl -f --user-unit openrwc.service # To view service logs\n"

echo -e "\nConfiguration file will be created in the user's home directory: ~/.config/OpenRWC/config.toml"
