# build/openrwc@-xfce.service
[Unit]
Description=Reddit Wallpaper Changer for GNU/Linux
Requires=display-manager.service
After=display-manager.service
StartLimitIntervalSec=0

[Service]
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/%i/.Xauthority"
Type=simple
ExecStart=/usr/bin/sh -c "DBUS_SESSION_BUS_ADDRESS=unix:path=$(find /run/user/ -maxdepth 1 -mindepth 1 | head -n 1)/bus /usr/bin/openrwc"
Restart=always
RestartSec=5
User=%i

[Install]
WantedBy=graphical.target
